# application.timer
[Unit]
Description=Timer for Web Service

[Timer]
OnUnitActiveSec=5m
Unit=polling_script.service

[Install]
WantedBy=timers.target
